# root/etc/s6-overlay/s6-rc.d/kavita/run
#!/usr/bin/with-contenv bash

cd /kavita || exit 1

cmd="./Kavita"
echo "*** executing => $cmd"
exec s6-setuidgid abc \
    $cmd
